# management.cron
# Company Cron

# Self-learning (every 5 mins)
*/5 * * * * cd /Users/macbookpro/.openclaw/workspace/company && python3 aggressive_learn.py >> /tmp/aggressive_learn.log 2>&1

# Operations (every 30 mins)
*/30 * * * * cd /Users/macbookpro/.openclaw/workspace/company && python3 auto_run.py >> /tmp/company.log 2>&1

# Daily report (8pm every day)
0 20 * * * /usr/bin/python3 /Users/macbookpro/.openclaw/workspace/company/weekly_telegram.py >> /tmp/daily_report.log 2>&1
