# management.cron
# Company Management Cron Jobs - Every 30 Minutes

# Management chain runs every 30 mins
*/30 * * * * /usr/bin/python3 /Users/macbookpro/.openclaw/workspace/company/management_chain.py daily >> /tmp/management.log 2>&1
