# management.cron
# Company Management Cron - PROFIT FOCUSED

# CEO → PM → Profit Generator every 30 minutes
*/30 * * * * /usr/bin/python3 /Users/macbookpro/.openclaw/workspace/company/management_chain.py daily >> /tmp/profit.log 2>&1

# Social media (Instagram/Telegram) every 6 hours
0 */6 * * * /usr/bin/python3 /Users/macbookpro/.openclaw/workspace/company/ig_poster.py >> /tmp/ig.log 2>&1
