# social.cron
0 */6 * * * /usr/bin/python3 /Users/macbookpro/.openclaw/workspace/company/social_poster.py >> /tmp/social_poster.log 2>&1
